# skills-stable/observability-contracts/references/custom-probe/probe.service
[Unit]
Description=<name> observability probe
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=/etc/<name>/env
WorkingDirectory=/opt/<name>
ExecStart=/opt/<name>/bin/<name>
User=<service-user>
Group=<service-user>
StateDirectory=<name>
LogsDirectory=<name>
RuntimeDirectory=<name>
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
